# etc/supervisord.d/swarm-visualizer.ini
[program:swarm-visualizer]
command = /opt/qnib/swarm-visualizer/bin/start.sh
directory = /opt/swarm-visualizer/
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
